# src/95-ka-lite-hotspot.rules
ACTION=="add", ENV{ID_VENDOR_ID}=="148f", RUN+="/usr/sbin/ka-lite-hotspot-service.sh start"
ACTION=="remove", ENV{ID_VENDOR_ID}=="148f", RUN+="/usr/sbin/ka-lite-hotspot-service.sh stop"