# src/ka-lite-hotspot.conf
description "WiFi Hotspot for KA Lite"

start on started gnome-session or started compiz
stop on shutdown

script
    /usr/sbin/ka-lite-hotspot --no-daemon start
end script
