# src/ka-lite-hotspot.conf
description "WiFi Hotspot for KA Lite"

start on started networking
start on stopping networking

script
    /usr/bin/ka-lite-hotspot --no-daemon start
end script
